# docker/dev/crontab
* * * * * php /var/www/html/cron/job.php indeed
