# docker/dev/crontab
* * * * * tail -f php /var/www/html/cron/job.php JobsController indeed
