# docker/dev/jobs.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * php /var/www/html/cron/job.php indeed > /var/www/html/cronErrors.txt
